# dists/freebsd/files/g15stats.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: g15stats
# REQUIRE: LOGIN g15daemon
# KEYWORD: shutdown
#
# Add these lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# g15stats_enable (bool):	Set to NO by default.
#				Set it to YES to enable doormand.

. /etc/rc.subr

name=g15stats
rcvar="${name}_enable"
load_rc_config $name

: ${g15stats_enable:="NO"}

command=%%PREFIX%%/bin/${name}
command_args='-d -i re0'

run_rc_command "$1"
